# scripts/crontab
*/2 * * * * /Users/lhl/Desktop/blueball/detector_a.py
* * * * * /Users/lhl/Desktop/blueball/scripts/post.a.sh
